# roles/docker/files/supervisord.conf
[supervisord]
nodaemon=true

[program:startxvfb]
command=/bin/sh startXvfb.sh
